# deploy/docker/supervisord/supervisord-celery_short.conf
[program:celery_short]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=10
numprocs=1
user=www-data
directory=/home/mediacms.io/mediacms
priority=400
startinorder=true
autorestart=true
stopasgroup=true
killasgroup=true
# Foreground worker so logs reach Docker stdout (was `celery multi start short1 short2 ... --logfile=...`
# which daemonized and wrote to a file inside the container, hiding all task/SMTP output).
# Concurrency 20 preserves the previous total of 2 nodes x 10 worker processes.
command=/home/mediacms.io/bin/celery -A cms worker -n short@%%h --loglevel=INFO --soft-time-limit=300 -c 20 -Q short_tasks
