# deploy/docker/supervisord/supervisord-celery_beat.conf
[program:celery_beat]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startsecs=10
numprocs=1
user=www-data
directory=/home/mediacms.io/mediacms
priority=300
startinorder=true
autorestart=true
stopasgroup=true
killasgroup=true
# Foreground beat so logs reach Docker stdout (was writing to /home/mediacms.io/mediacms/logs/celery_beat.log).
command=/home/mediacms.io/bin/celery -A cms beat --loglevel=INFO
